# modelos/hermes-sdr-watchdog.timer
[Unit]
Description=Verifica a cada minuto se o webhook do Hermes SDR esta no ar

[Timer]
# Começa 90s após o boot (dá tempo do túnel estabelecer a borda antes da 1ª
# checagem, evitando um restart desnecessário logo no boot) e repete a cada 60s.
OnBootSec=90
OnUnitActiveSec=60
AccuracySec=10s
Unit=hermes-sdr-watchdog.service

[Install]
WantedBy=timers.target
